# cronjob-mdadm.service
[Unit]
Description=mdadm-checkarray cronjob

[Service]
ExecStart=/sbin/mdadm-checkarray --cron --all --quiet
